# charms/trusty/limeds-livinglabs/templates/upstart.conf
#!upstart
description "{{description}}"

start on startup
stop on shutdown

respawn

script
{{command}}
end script
