# charms/trusty/rest2kafka/templates/upstart.conf
#!upstart
description "{{description}}"

start on startup
stop on shutdown

respawn

script
su - {{user}} -c 'export DEBUG={{debug}}; {{command}}'
end script
