# charms/layers/jupyter-notebook/templates/upstart.conf
#!upstart
description "jupyter"

start on startup
stop on shutdown

respawn

script
jupyter notebook --no-browser --config /opt/jupyter/jupyter_notebook_config.py
end script
