# charms/layers/jupyter-notebook/templates/upstart.conf
#!upstart
description "jupyter"

start on startup
stop on shutdown

respawn

script
jupyter notebook --no-browser
end script
