# charms/trusty/wetty/templates/upstart.conf
#!upstart
description "snappass"

start on startup
stop on shutdown

setuid ubuntu
setgid ubuntu

respawn

script
cd /opt/wetty/
node app.js --sslkey key.pem --sslcert cert.pem -p 3000
end script
