# charms/trusty/apache-kafka/templates/upstart.conf
#!upstart
description "kafka"

start on startup
stop on shutdown

setuid kafka
setgid hadoop

respawn

script
"{{kafka_bin}}/kafka-server-start.sh" "{{kafka_conf}}/server.properties"
end script
