# charms/trusty/http-test/templates/upstart.conf
#!upstart
description "http-test"

start on startup
stop on shutdown

respawn

script
export DEBUG=False
/opt/http-test/http-test.py
end script
