# charms/layers/tengu-api/templates/upstart.conf
#!upstart
description "tokin"

start on startup
stop on shutdown

respawn

script
su - ubuntu -c '/opt/tengu-api/tengu-api.py'
end script
