# charms/trusty/hauchiwa/templates/upstart.conf
#!upstart
description "h_api"

start on startup
stop on shutdown

respawn

script
su - {{user}} -c '{{tengu_dir}}/scripts/h_api.py'
end script
