# charms/trusty/apache-flume-kafka/templates/upstart.conf
#!upstart
description "flume"

start on startup
stop on shutdown

setuid flume

respawn

script
"{{flume}}/bin/flume-ng" agent -c "{{flume_conf}}" -f "{{flume_conf}}/flume.conf" -n a1
end script
