# charms/trusty/spark-notebook/templates/upstart.conf
#!upstart
description "spark-notebook"

start on startup
stop on shutdown


respawn

script
su - ubuntu -c '/opt/spark-notebook/bin/spark-notebook'
end script
