# charms/trusty/spark-notebook/templates/upstart.conf
#!upstart
description "snappass"

start on startup
stop on shutdown

setuid ubuntu
setgid ubuntu

respawn

script
snappass
end script
