# charms/layers/rest2maas/templates/upstart.conf
#!upstart
description "rest2maas"

start on startup
stop on shutdown

respawn

script
export DEBUG=False
/opt/rest2maas/rest2maas.py
end script
