# charms/trusty/tokin/templates/upstart.conf
#!upstart
description "tokin"

start on startup
stop on shutdown

respawn

script
export DEBUG=False
su - ubuntu -c '/opt/tokin/tokin.py'
end script
