# charms/layers/tokin/templates/upstart.conf
#!upstart
description "tokin"

start on startup
stop on shutdown

setuid ubuntu
setgid ubuntu

respawn

script
/opt/tokin/tokin.py
end script
